# adminbot-watcher.path
[Unit]
Description=Watch admin bot code for changes

[Path]
PathModified=~/telegram-claude-bot/admin_bot/commands.py
PathModified=~/telegram-claude-bot/admin_bot/bridge.py
PathModified=~/telegram-claude-bot/admin_bot/cognitive.py
PathModified=~/telegram-claude-bot/admin_bot/chat.py
PathModified=~/telegram-claude-bot/admin_bot/config.py
PathModified=~/telegram-claude-bot/admin_bot/domains.py
PathModified=~/telegram-claude-bot/admin_bot/__main__.py
PathModified=~/telegram-claude-bot/admin_bot.py
Unit=adminbot-restart.service

[Install]
WantedBy=default.target
