# personabots-watcher.path
[Unit]
Description=Watch persona bot code for changes

[Path]
PathModified=~/telegram-claude-bot/bot_base.py
PathModified=~/telegram-claude-bot/memory.py
PathModified=~/telegram-claude-bot/conversation_compressor.py
Unit=personabots-restart.service

[Install]
WantedBy=default.target
